# Ex1/create_file_rule.rules
ACTION=="add",KERNEL=="sd[a-z]",RUN+="/home/nicolas/create_file_rule.sh"
